# 08 Linux/custom/board/rootfs_overlay/etc/profile.d/prompt.sh
color1="\[\033[1;34m\]"
color2="\[\033[0m\]"
export PS1="${color1}\u@\h${color2}:\w\$ "
unset color1 color2
