# 08 Linux/debian/overlay/etc/systemd/system/my-keyboard-layout.service
[Unit]
Description=Switch Keyboard Layout

[Service]
Type=oneshot
ExecStart=/usr/bin/loadkeys de-latin1

[Install]
WantedBy=default.target

