# 08 Linux/debian/overlay/etc/systemd/system/my-regenerate-ssh-keys.service
[Unit]
Description=Regenerate SSH Keys

[Service]
Type=oneshot
ExecStart=/usr/bin/ssh-keygen -A

[Install]
WantedBy=default.target

